# atm9-restart.service
[Unit]
Description=Triggers ATM9 Restart Script

[Service]
Type=oneshot
User=ascentius
Group=ascentius
ExecStart=/var/atm9/restart_script.sh
